# rtmp-server/rtmp-manager.service
[Unit]
Description=RTMP Restream Manager
After=network.target nginx.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/rtmp-manager
ExecStart=/usr/bin/node /opt/rtmp-manager/server.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=3000
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
